# deploy/systemd/climate-controller.service
[Unit]
Description=Sagrada Climate Heating Controller
After=network.target mariadb.service mosquitto.service
Wants=mariadb.service mosquitto.service

[Service]
Type=simple
User=peterbutler
WorkingDirectory=/home/peterbutler/sagrada
Environment=CLIMATE_ENV=sagrada
ExecStart=/home/peterbutler/sagrada/venv/bin/python scripts/run_controller.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
