# deploy/systemd/readings-aggregator.service
[Unit]
Description=Sagrada Readings Aggregator
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=peterbutler
WorkingDirectory=/home/peterbutler/sagrada
Environment=CLIMATE_ENV=sagrada
ExecStart=/home/peterbutler/sagrada/venv/bin/python scripts/run_aggregator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
