# deploy/systemd/sagrada-api.service
[Unit]
Description=Sagrada TypeScript API Server
After=network.target mariadb.service mosquitto.service
Wants=mariadb.service mosquitto.service

[Service]
Type=simple
User=peterbutler
WorkingDirectory=/home/peterbutler/sagrada/api
EnvironmentFile=/home/peterbutler/sagrada/config/.env
Environment=NODE_ENV=production
Environment=PATH=/home/peterbutler/.nvm/versions/node/v22.19.0/bin:/usr/bin:/bin
ExecStart=/home/peterbutler/.nvm/versions/node/v22.19.0/bin/npx tsx src/index.ts
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
